# etc/supervisor/language_tool.conf
[program:language_tool]
command = java -cp languagetool-server.jar org.languagetool.server.HTTPServer --port 7001
directory = /home/jmas/languagetool
user = jmas
