# cfg/supervisor/recursos.conf
[program:recursos]
command = /home/jmas/web/python-env/bin/gunicorn web_search:app -b localhost:8000
directory = /home/jmas/web/recursos
user = jmas
environment=PATH="/home/jmas/web/python-env/bin/"
